# packaging/keyzer-tilt-shim.service
[Unit]
Description=KEYZER Naga Pro wheel-tilt shim (hidraw -> uinput)
Documentation=https://git.jagbhandal.com/jagdeep.bhandal/keyzer
# Best-effort: the shim self-detects the device and exits cleanly if absent.

[Service]
Type=simple
# @KEYZER_DIR@ is replaced with your KEYZER checkout path by install.sh.
ExecStart=/usr/bin/python3 @KEYZER_DIR@/app/tilt_shim.py
Restart=on-failure
RestartSec=3
# It only reads the Razer mouse hidraw and writes its own uinput device.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=default.target
